# scripts/systemd/torrentfs.service
[Unit]
Description=TorrentFS daemon (user)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p %h/torrentfs/torrents %h/torrentfs/cache
ExecStart=/usr/bin/torrentfsd --torrent-dir %h/torrentfs/torrents --cache %h/torrentfs/cache --socket %t/torrentfsd.sock
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
